# trunk/installers/irobot/SkypeSparkyServer/leopard/Package_Root/Sparky/sparky.cron
#check if control.php is running
#
* * * * * php /Sparky/shell_check_sparky_irobot.php &> /tmp/shell_check_sparky_irobot.php.log

